# bin/jutty.service
# systemd unit file
#
# place in /etc/systemd/system
# systemctl enable jutty.service
# systemctl start jutty.service

[Unit]
Description=Web-based ssh/telnet client
After=network.target

[Service]
User=root
Group=root

# change WorkingDirectory to directory with project
WorkingDirectory=/home/kj/PhpstormProjects/jutty
ExecStart=/usr/bin/node dist/app.js

Restart=always
RestartSec=10

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=nodejs-jutty

[Install]
WantedBy=multi-user.target
